# setup/carbon-tracker.service
[Unit]
Description=Carbon Tracker Web (Gunicorn)
After=network.target

[Service]
User=ec2-user
WorkingDirectory=/opt/carbon-tracker
EnvironmentFile=/opt/carbon-tracker/.env
ExecStart=/opt/carbon-tracker/run_gunicorn.sh
Restart=always

[Install]
WantedBy=multi-user.target